# photobot_src/supervisord_conf/tunnels.conf
[program:tunneler]
command=loclx tunnel --config /var/photobot/config/tunnels.yml
autorestart=true